# ops/systemd/paywifi-watchdog.timer
[Unit]
Description=Run PAYWIFI watchdog every 60s

[Timer]
OnBootSec=60
OnUnitActiveSec=60
AccuracySec=5s
Unit=paywifi-watchdog.service

[Install]
WantedBy=timers.target
