# ops/systemd/paywifi-watchdog.service
[Unit]
Description=PAYWIFI infrastructure watchdog (self-heal IFB/tc/nft/services)
After=paywifi-shape.service network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/paywifi-watchdog.sh
